# sssd/oddjobd.service
.include /usr/lib/systemd/system/oddjobd.service
[Service]
ExecStartPre=/bin/systemctl start dbus.service
