# Bind/cfg_files/supervisord.d/sshd.ini
[program:sshd]
command=/etc/init.d/sshd
